# etc/udev/rules.d/99-mxc_asrc.rules
KERNEL=="mxc_asrc", SUBSYSTEM=="misc", ACTION=="add", RUN+="/bin/chmod 777 /dev/mxc_asrc"
